# docker/prod/app/supervisor/app.conf
[supervisord]
nodaemon=true

[program:app]
<<<<<<< HEAD
command=bundle exec pumactl start
=======
command=bundle exec puma -C config/puma.rb
>>>>>>> 3fbd0abd95d21a8f62e12ce666145f8f00c17160
autostart=true
autorestart=true
stopsignal=TERM
user=root
directory=/var/www/mtm/
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stopsignal=TERM
user=root
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:sshd]
command=/usr/local/bin/docker-entrypoint.sh
